# deploy/iggy-server.service
# Apache Iggy server — llmbus broker (no Docker; §9b).
#
# The binary is BUILT FROM SOURCE at tag server-0.8.0 (no prebuilt 0.8.0 binary
# exists) and installed to /usr/local/bin/iggy-server — see deploy/README.md.
# Config comes entirely from env (deploy/iggy.env); IGGY_* overrides Iggy's built-in
# defaults, so no config.toml is needed.
#
# Install:  sudo cp deploy/iggy-server.service /etc/systemd/system/
#           sudo systemctl daemon-reload && sudo systemctl enable --now iggy-server
[Unit]
Description=Apache Iggy server (llmbus broker)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=bartek
Group=bartek
# Data lives at IGGY_SYSTEM_PATH (=/var/lib/iggy, set in the env file); WorkingDirectory
# just catches any relative writes. Create it first: sudo mkdir -p /var/lib/iggy &&
# sudo chown bartek:bartek /var/lib/iggy
WorkingDirectory=/var/lib/iggy
EnvironmentFile=/home/bartek/Projects/llmbus/deploy/iggy.env
ExecStart=/usr/local/bin/iggy-server
# Iggy is a thread-per-core io_uring server: it raises scheduling priority (hwloc)
# and locks memory. Grant the capability + memlock; do NOT add a SystemCallFilter —
# it would block the io_uring syscalls and the server won't start.
AmbientCapabilities=CAP_SYS_NICE
LimitMEMLOCK=infinity
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
